# systemd-services/claude-mcp.service
[Unit]
Description=Claude Code MCP SSE Server
After=network.target

[Service]
Type=simple
User=mcp-user
WorkingDirectory=/opt/stacks/claude-code-mcp
# Load environment from secure file
EnvironmentFile=/etc/claude-mcp/environment
ExecStart=/usr/bin/node /opt/stacks/claude-code-mcp/dist/server.js --http
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target